# systemd/freepbx-nightly-scrape.timer
[Unit]
Description=FreePBX Nightly Scrape (daily at 2:00 AM)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
