# systemd/freepbx-tools.service
[Unit]
Description=FreePBX Tools Web Services
After=network-online.target apache2.service
Wants=network-online.target
StartLimitIntervalSec=180
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/var/www/freePBX_Tools
ExecStart=/var/www/freePBX_Tools/scripts/start_services.sh
ExecStop=/usr/bin/python3 /var/www/freePBX_Tools/scripts/stop_all_web_apps.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freepbx-tools
TimeoutStartSec=240

[Install]
WantedBy=multi-user.target
